# resources/hbase/supervisor-regionserver.conf
[program:regionserver]
command=/opt/hbase/bin/hbase regionserver start
directory=/opt/hbase
stdout_logfile=/var/log/hbase/regionserver-stdout.log
stderr_logfile=/var/log/hbase/regionserver-stderr.log